# clients/c8y-devicecert-renewer/services/systemd/c8y-devicecert-renewer.target
[Unit]
Description=Target to Cumulocity self-signed certificate renewer
